# supervisord.conf
[supervisord]
nodaemon=true

[program:nodescreenservice]
command=node /src/app.js
stderr_logfile = /servicelogs/nodescreenservice-stderr.log
stdout_logfile = /servicelogs/nodescreenservice-stdout.log